# Library/LaunchAgents/net.devpi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>net.devpi</string>

        <key>EnvironmentVariables</key>
        <dict>
            <key>PYTHONPATH</key>
            <string></string>
        </dict>

        <key>ProgramArguments</key>
        <array>
            <string>/data/VENV/sax/bin/python</string>
            <string>/data/VENV/sax/bin/devpi-server</string>

            <string>--host</string>
            <string>0.0.0.0</string>

            <string>--port</string>
            <string>3141</string>
            <string>--serverdir</string>
            <string>/data/devpi_root</string>

        </array>

        <key>RunAtLoad</key>
        <true/>
    </dict>
</plist>
